# launchd/com.contentengine.x.midday.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.contentengine.x.midday</string>
  <key>WorkingDirectory</key>
  <string>/Users/admin/content-engine</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/admin/content-engine/code/posting/run-x-slot.js</string>
    <string>midday</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>12</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/admin/content-engine/output/logs/x-midday.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/admin/content-engine/output/logs/x-midday.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
